# web_server/runz-web.sh.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:          monoserve.sh
# Required-Start:    $local_fs $syslog $remote_fs
# Required-Stop:     $local_fs $syslog $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start fastcgi mono server with hosts
### END INIT INFO
<% @mono_ver =  ENV['MONO_VERSION']%>
MONO_PREFIX=/opt/mono-<%= @mono_ver %>
export DYLD_LIBRARY_PATH=$MONO_PREFIX/lib:$DYLD_LIBRARY_PATH
export LD_LIBRARY_PATH=$MONO_PREFIX/lib:$LD_LIBRARY_PATH
export C_INCLUDE_PATH=$MONO_PREFIX/include
export ACLOCAL_PATH=$MONO_PREFIX/share/aclocal
export PKG_CONFIG_PATH=$MONO_PREFIX/lib/pkgconfig
export MONO_IOMAP=case

PATH=/opt/mono-<%= @mono_ver %>/frog/bin:/opt/mono-<%= @mono_ver %>/bin:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
NAME=monoserver
DESC=monoserver

MONOSERVER=$(which xsp4)
MONOSERVER_PID=$(ps auxf | grep xsp4.exe | grep -v grep | awk '{print $2}')

WEBAPPS="/stat:/home/web/runz/web"

case "$1" in
        start)
                if [ -z "${MONOSERVER_PID}" ]; then
                        echo "starting mono server"
                        MONO_OPTIONS="--debug " ${MONOSERVER} --applications ${WEBAPPS} --nonstop &
                        echo "mono server started"
                else
                        echo ${WEBAPPS}
                        echo "mono server is running"
                fi
        ;;
        stop)
                if [ -n "${MONOSERVER_PID}" ]; then
                        kill ${MONOSERVER_PID}
                        echo "mono server stopped"
                else
                        echo "mono server is not running"
                fi
        ;;
esac

exit 0
